# vendor/github.com/clearcontainers/proxy/cc-proxy.service.in
[Unit]
Description=Clear Containers Proxy
Documentation=https://github.com/clearcontainers/proxy

[Service]
ExecStart=@libexecdir@/cc-proxy

[Install]
WantedBy=multi-user.target
