# vendor/github.com/clearcontainers/proxy/cc-proxy.service.in
[Unit]
Description=Clear Containers Proxy
Documentation=https://github.com/clearcontainers/proxy

[Service]
ExecStart=@libexecdir@/clear-containers/cc-proxy
LimitNOFILE=infinity

[Install]
WantedBy=multi-user.target
