# vendor/github.com/clearcontainers/proxy/cc-proxy.socket.in
[Unit]
Description=Clear Containers Proxy Socket
Documentation=https://github.com/clearcontainers/proxy
PartOf=cc-proxy.service

[Socket]
ListenStream=@localstatedir@/run/clear-containers/proxy.sock
DirectoryMode=0770
SocketMode=0660

[Install]
WantedBy=sockets.target
