# vendor/github.com/01org/ciao/networking/ciao-cnci-agent/scripts/ciao-cnci-agent.service
[Unit]
Description=Ciao CNCI Agent
After=network.target

[Service]
ExecStartPre=-/usr/bin/ethtool -L enp0s2 combined 4
ExecStart=/usr/sbin/ciao-cnci-agent -server auto -v 3
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartSec=31s

[Install]
WantedBy=default.target
