# .ci/upstart-services/cc-proxy.conf
description "Clear Containers Proxy"

# start in normal runlevels when disks are mounted and networking is available
start on runlevel [2345]

# stop on shutdown/halt, single-user mode and reboot
stop on runlevel [016]

exec /usr/libexec/clearcontainers/cc-proxy -v 3
