# vendor/github.com/01org/ciao/_DeploymentAndDistroPackaging/ansible/roles/keystone/templates/docker-keystone.service.j2
[Unit]
Description=Keystone Container for CIAO
Requires={{ docker_service }}
After={{ docker_service }}

[Service]
Restart=always
ExecStart=/usr/bin/docker start -a ciao-keystone
ExecStop=/usr/bin/docker stop -t 2 ciao-keystone

[Install]
WantedBy=default.target
