# vendor/github.com/01org/cc-oci-runtime/data/cc-agent.service
[Unit]
Description=Clear Container Agent

[Service]
StandardOutput=tty
Type=simple
ExecStart=/bin/hyperstart
ExecStop=/usr/bin/systemctl --force poweroff
FailureAction=poweroff
