# vendor/github.com/01org/cc-oci-runtime/data/container.target
[Unit]
Description=Containers Target
Requires=basic.target
Conflicts=rescue.service rescue.target
After=basic.target rescue.service rescue.target
AllowIsolate=yes
Wants=systemd-logind.service
Wants=container-workload.service
Wants=systemd-networkd.service
Wants=systemd-resolved.service
Wants=opt-rootfs.mount
Wants=opt-rootfs-proc.mount
Wants=opt-rootfs-sys.mount
Wants=hyperstart.service
