# vendor/github.com/01org/cc-oci-runtime/data/container-workload.service
[Unit]
Description=Container Workload
After=opt-rootfs-resolv.service
After=opt-rootfs-proc.mount
After=opt-rootfs-sys.mount

[Service]
RootDirectory=/opt/rootfs
RootDirectoryStartOnly=yes
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
EnvironmentFile=-/opt/rootfs/.containerenv
StandardInput=tty
StandardOutput=tty
PrivateDevices=yes
Type=simple
ExecStart=/.containerexec
ExecStop=/usr/bin/systemctl --force poweroff
FailureAction=poweroff
