# vendor/github.com/01org/cc-oci-runtime/data/opt-rootfs-proc.mount
[Unit]
After=opt-rootfs.mount

[Mount]
What=proc
Where=/opt/rootfs/proc
Type=proc
