# vendor/github.com/01org/cc-oci-runtime/data/opt-rootfs-sys.mount
[Unit]
After=opt-rootfs.mount

[Mount]
What=sys
Where=/opt/rootfs/sys
Type=sysfs
Options=nosuid,noexec,nodev,ro
