# vendor/github.com/01org/cc-oci-runtime/data/opt-rootfs.mount
[Mount]
What=rootfs
Where=/opt/rootfs
Type=9p
Options=trans=virtio,version=9p2000.L,msize=131072
