# vendor/github.com/01org/ciao/_DeploymentAndDistroPackaging/ansible/roles/ciao-controller/templates/ciao-controller.service.j2
[Unit]
Description=CIAO Controller Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ciao-controller --cacert=/etc/pki/ciao/CAcert-{{ ciao_controller_fqdn }}.pem \
                                   --cert=/etc/pki/ciao/cert-Controller-localhost.pem \
                                   --v 3
User=ciao
Restart=always

[Install]
WantedBy=multi-user.target
