# vendor/github.com/01org/cc-oci-runtime/proxy/cc-proxy.service.in
[Unit]
Description=Clear Containers Proxy
Documentation=https://github.com/01org/cc-oci-runtime/proxy

[Service]
ExecStart=@libexecdir@/cc-proxy

[Install]
WantedBy=multi-user.target
