# vendor/github.com/01org/cc-oci-runtime/proxy/cc-proxy.socket.in
[Unit]
Description=Clear Containers Proxy Socket
Documentation=https://github.com/01org/cc-oci-runtime/proxy
PartOf=cc-proxy.service

[Socket]
ListenStream=@localstatedir@/run/cc-oci-runtime/proxy.sock

[Install]
WantedBy=sockets.target
